# 1.1-udev-rule/71-usb.rules
# match serial with neccessary
ACTION=="add", ATTRS{serial}=="M2YGSJ81", RUN+="/usr/bin/copy-logs.sh"